# ops/com.codex.stock-daily-report.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.codex.stock-daily-report</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/huangchenyao/Documents/Codex/2026-04-25/top10-3-5-10-20-top10/scripts/run_daily.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>21</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/huangchenyao/Documents/Codex/2026-04-25/top10-3-5-10-20-top10/output/launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/huangchenyao/Documents/Codex/2026-04-25/top10-3-5-10-20-top10/output/launchd.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/huangchenyao/Documents/Codex/2026-04-25/top10-3-5-10-20-top10</string>
    <key>RunAtLoad</key>
    <false/>
  </dict>
</plist>
